# systemd/devicectl-stop@.service
[Unit]
Description=Device Control Stop %I
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/devicectl %I stop
